# dotfiles/systemd/awsso-agent.service
[Unit]
Description=Awsso agent server

[Service]
ExecStart=/home/psanford/bin/awsso server

[Install]
WantedBy=default.target
