# dotfiles/systemd/brightness-buddy.service
[Unit]
Description=Brightness budy

[Service]
ExecStart=/home/psanford/bin/brightness-buddy

[Install]
WantedBy=default.target
